# config/profile.d/man.sh
# Man pager
# Define how to display man pages
# export MANPAGER="less -sMJ +Gg"
export MANPAGER="sh -c 'col -bx | bat -l man -p'" # batcat
