# config/profile/profile.d/hist.sh
# Configure how the history is written
export HISTSIZE=10000
export HISTFILESIZE=100000
export SAVEHIST=10000
export HISTCONTROL=ignoreboth
export HISTIGNORE='ls:ll:la:lc:lt:cd:gt:up:todo:exit:clear:cls:hist:history'
