# config/profile/profile.d/ripgrep.sh
#!/usr/bin/env sh
# Ripgrep Configuration: Set the path to the Ripgrep configuration file.

# Define the location of the Ripgrep configuration file
export RIPGREP_CONFIG_PATH="$XDG_CONFIG_HOME/ripgrep/ripgreprc"
