# config/profile/profile.d/zk.sh
#!/usr/bin/env sh
# zk: Plain text note-taking assistant
# https://github.com/zk-org/zk

# Set the notebook directory for zk notes.
export ZK_NOTEBOOK_DIR="$HOME/Documents/notes"
